# misc/upstart.conf
description "fermpi"

start on runlevel [2345]
stop on runlevel [6]

respawn

script
  export NODE_ENV=production
  cd /home/pi/fermpi
  exec npm start
end script
